# vendor/init.d/social-pay-rpc
#!/bin/sh
#
# go - this script starts and stops the golang bin daemon
#
# chkconfig:   - 85 15 
# description:  Golang Echo Serveris an HTTP(S) server
# processname: $1
 
# Source function library.
. /etc/rc.d/init.d/functions
 
# Source networking configuration.
. /etc/sysconfig/network
 
# Check that networking is up.
[ "$NETWORKING" = "no" ] && exit 0
 
path="/data/www/social_backend_pay/bin/social-pay-rpc-server"
prog=$(basename $path)
 
start() {
    [ -x $path ] || exit 5
    echo -n $"Starting $prog: "
    cmd=`nohup $path >> /data/logs/$prog.log 2>&1 &`
    $cmd && success || failure $"$prog start"
    retval=$?
    echo
    return $retval
}
 
stop() {
    echo -n $"Stopping $prog: "
    killproc $prog -9
    retval=$?
    echo
    return $retval
}
 
restart() {
    stop
    start
}
 
reload() {
    echo -n $"Reloading $prog: "
    killproc $prog -USR2
    retval=$?
    echo
    return $retval
}

rh_status() {
    status $prog
}

rh_status_q() {
    rh_status >/dev/null 2>&1
}
 
case "$1" in
    start)
        rh_status_q && exit 0
        $1
        ;;
    stop)
        rh_status_q || exit 0
        $1
        ;;
    restart)
        $1
        ;;
    reload)
        rh_status_q || exit 7
        $1
        ;;
    *)
        echo $"Usage: $0 progname {start|stop|reload|restart}"
        exit 2
esac
